# playbook/files/socat-upstart.conf
# socket forwarding for docker
start on runlevel [2345]

respawn

exec socat TCP-LISTEN:2375,fork,reuseaddr UNIX:/var/run/docker.sock